# scripts/hdparm.service
[Unit]
Description=hdparm sleep

[Service]
Type=oneshot
ExecStart=/sbin/hdparm -q -S 120 -Y /dev/sda

[Install]
WantedBy=multi-user.target
